# container-files/etc/services.d/apache/run
#!/bin/execlineb -P
#/usr/sbin/httpd -DFOREGROUND -c "SetEnv MYSQL_USER ${MYSQL_USER_ENV}" -c "SetEnv MYSQL_ROOT_PASSWORD ${MYSQL_ROOT_PASSWORD_ENV}" -c "SetEnv MYSQL_DATABASE ${MYSQL_DATABASE_ENV}" -c "SetEnv MYSQL_PASSWORD ${MYSQL_PASSWORD_ENV}" -c "SetEnv MYSQL_DATABASE ${MYSQL_DATABASE_ENV}" -c "SetEnv DATABASE_PREFIX ${DATABASE_PREFIX_ENV}"