# container-files/etc/services.d/apache/run
#!/bin/execlineb -P
/usr/sbin/httpd -DFOREGROUND
