# container-files/etc/services.d/apache/run
#!/bin/execlineb -P
/usr/sbin/httpd -DFOREGROUND -D MYSQL_USER=MYSQL_ROOT_PASSWORD
